# irq_service/irqservice-upstart.conf
description "IRQ web service"
author "matthew.machaj@gmail.com"

start on runlevel [2345]
stop on runlevel [016]

exec irqservice

